# .config/systemd/user/hyprland.service
[Unit]
Description=Hyprland

[Service]
ExecStart=Hyprland

[Install]
WantedBy=default.target